# emacs.service
[Service]
User=root
ExecStart="emacsclient -c"

[Install]
WantedBy=multi-user.target
